# systemd/cc-switch@.service
[Unit]
Description=CC-Switch Server
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/opt/opencode
ExecStart=/opt/opencode/cc-switch-main/src-tauri/target/release/cc-switch-server
Restart=always
RestartSec=10
StandardOutput=append:/var/log/opencode/cc-switch.log
StandardError=append:/var/log/opencode/cc-switch.log

[Install]
WantedBy=multi-user.target
